# task-relay/deploy/systemd/task-relay-db-check.service
[Unit]
Description=task-relay database integrity check
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=oneshot
User=task-relay
EnvironmentFile=/etc/task-relay/task-relay.env
WorkingDirectory=/var/lib/task-relay
ExecStart=/usr/local/bin/task-relay db-check
